# image-files/supervisord-tomcat.conf
[program:tomcat]
command=/start-tomcat.sh
numprocs=1
autostart=true
autorestart=true
priority=20